# audio/99-hdmi-audio.rules
ACTION=="add|change", SUBSYSTEM=="drm", RUN+="/usr/local/bin/audio-switch.sh"